# distributed-conf/distributed-supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:initsh]
command=bash -c '/root/conf/distributed-init.sh'
